# redis-memcached/loaddata.service
[Unit]
Description=Load Dummy Data into Redis

[Service]
Type=oneshot
ExecStart=redis-memcached\script.sh
