# redis-memcached/loaddata.timer
[Unit]
Description=Load Dummy Data into Redis Timer

[Timer]
OnUnitActiveSec=5m
Unit=load_data.service

[Install]
WantedBy=multi-user.target
